# upstart/mosquitto.conf
#!upstart
# Ubuntu upstart file at /etc/init/mosquitto.conf

start on startup
stop on shutdown

exec mosquitto -c /etc/mosquitto/mosquitto.conf
