# upstart/mosquitto.conf
#!upstart
# Ubuntu upstart file at /etc/init/mosquitto.conf

#respawn
#respawn limit 15 5

start on startup
stop on shutdown

exec mosquitto -c /etc/mosquitto/mosquitto.conf
