# upstart/mosquitto.conf
#!upstart
# Ubuntu upstart file at /etc/init/mosquitto.conf

start on (net-device-up IFACE=ppp0)
stop on shutdown

exec mosquitto -c /etc/mosquitto/mosquitto.conf
